# supervisor.conf
[program:temperature-logger]
process_name=temperature-logger
directory=/home/pi/temperature-logger
environment=PATH="/home/pi/temperature-logger/bin"
command=/home/pi/temperature-logger/bin/python3 /home/pi/temperature-logger/temperature-logger.py
startsecs=5
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/temperature-logger-stderr.log
stdout_logfile=/var/log/supervisor/temperature-logger-stdout.log